# docker/supervisor/prod.conf
[include]
files = /etc/supervisor/common.conf

[program:octane]
command=php /var/www/html/artisan octane:frankenphp
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
